# docker/etc/s6/vikunja/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/vikunja > /dev/null
    exec su-exec $USER /app/vikunja/vikunja
popd
