# docker/etc/services.d/vikunja/run
#!/bin/execlineb -P
/app/vikunja/vikunja
